# php/supervisor/server.ini
[program: nginx]
command=/usr/sbin/nginx  -g 'daemon off;'
autorestart=true ;
autostart=true ;
stderr_logfile=/data/supervisor/logs/nginx.error.log ;
stdout_logfile=/data/supervisor/logs/nginx.stdout.log ;
user=root ;
stopsignal=INT
startsecs=10 ;
startretries=5 ;
stopasgroup=true

[program: php-fpm]
command=/usr/sbin/php-fpm
process_name=%(program_name)s ;
autorestart=true ;
autostart=true ;
stderr_logfile=/data/supervisor/logs/php.error.log ;
stdout_logfile=/var/logs/supervisor/php.stdout.log ;
user=root ;
stopsignal=INT
startsecs=10 ;
startretries=5 ;
stopasgroup=true