# php/supervisor/server.ini
[program: nginx]
command=/usr/sbin/nginx  -g 'daemon off;'
autorestart=true ;
autostart=true ;
stderr_logfile=/var/logs/supervisor/nginx.error.log ;
stdout_logfile=/var/logs/supervisor/nginx.stdout.log ;
user=root ;
stopsignal=INT
startretries=5 ;
stopasgroup=true

[program: php-fpm]
command=bash -c "/usr/sbin/php-fpm --fpm-config /etc/php-fpm.conf --pid /var/run/php-fpm.pid"
process_name=%(program_name)s ;
autostart=true
autorestart=true
startretries=5
exitcodes=0,2,70
stopsignal=QUIT
stopwaitsecs=5
stderr_logfile=/var/logs/supervisor/php.error.log 
stdout_logfile=/var/logs/supervisor/php.stdout.log 